# rootfs/build/services/alist.conf
[program:alist]

command = "/build/bin/start_alist"
autostart = true
autorestart = true