# rootfs/build/services/alist.conf
[program:alist]

command = "/build/bin/start_alist"
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /dev/stdout
