# rootfs/build/services/wireguard.conf
[program:wireguard]

command = "/build/bin/start_wireguard"
startsecs = 0
autorestart = false
startretries = 1
redirect_stderr = true
stdout_logfile = /dev/stdout
